# deploy/systemd/dnp3_master_io_v2.service
[Unit]
Description=DNP3 Master Fleet Monitor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=cic
WorkingDirectory=/home/cic/projects/dnp3/dnp3_master_io_v2
Environment=MASTER_LOG=info
ExecStart=/home/cic/projects/dnp3/dnp3_master_io_v2/target/release/dnp3_master_io_v2
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
